# etc/systemd/system/daily-backup-cleanup.service
[Unit]
Description=Daily Backup Cleanup

[Service]
Type=simple
Nice=10
ExecStart=/usr/local/bin/backup cleanup
Environment="HOME=/root"
