# etc/systemd/system/daily-backup-mysql.service
[Unit]
Description=Daily MySQL Backup

[Service]
Type=simple
Nice=10
ExecStart=/usr/local/bin/backup mysql
Environment="HOME=/root"
