# etc/systemd/system/daily-backup.service
[Unit]
Description=Daily Backup

[Service]
Type=simple
Nice=10
ExecStart=/usr/local/bin/backup run
Environment="HOME=/root"
